# debian/pipaos-tools.postinst
#!/bin/bash
#
# debian postinst script
#

case "$1" in

    configure)

        # Activate the bootup loog
        ln -sfv -r /usr/share/pipaos-tools/logos/linux-pipaos-800x600.png /usr/share/pipaos-tools/logos/default.png
        systemctl enable pipaos-bootlogo.service
        ;;
esac

#DEBHELPER#

exit 0
